# cmd/cri-resmgr/cri-resource-manager.service.in
[Unit]
Description=A CRI proxy with (hardware) resource aware container placement policies.
Documentation=https://github.com/intel/cri-resource-manager
Before=kubelet.service

[Service]
Type=simple
EnvironmentFile=__DEFAULTDIR__/cri-resource-manager
ExecStart=__BINDIR__/cri-resmgr $CONFIG_OPTIONS $POLICY_OPTIONS
Restart=always

[Install]
WantedBy=multi-user.target
